# .claude/worktrees/bold-ardinghelli/public/data/briefings/TB-2026-03-28-001.json
{
  "briefing_id": "TB-2026-03-28-001",
  "title": "Aqua Security Trivy Supply Chain Attack CVE-2026-33634 (CVSS 9.4): TeamPCP Injects CanisterWorm Credential-Stealing Malware into 76 Version Tags — Blockchain-Based C2 Infrastructure Propagates to 141+ npm Packages",
  "generated_date": "2026-03-28",
  "status": "ACTIVE",
  "severity": {
    "level": "Critical",
    "cvss": 9.4,
    "cisa_kev": true,
    "cisa_deadline": "2026-04-08"
  },
  "cves": [
    "CVE-2026-33634"
  ],
  "tags": [
    "supply-chain",
    "trivy",
    "aquasecurity",
    "devops",
    "ci-cd",
    "credential-theft",
    "worm",
    "blockchain-c2",
    "container-security",
    "github-actions"
  ],
  "affected_products": [
    "Aqua Security Trivy v0.69.4",
    "aquasecurity/trivy-action v0.0.1 through v0.34.2",
    "aquasecurity/setup-trivy pre-v0.2.6",
    "All downstream npm packages infected by CanisterWorm"
  ],
  "sources_analyzed": 15,
  "sources": [
    {
      "name": "Tenable CVE-2026-33634 Entry",
      "url": "https://www.tenable.com/cve/CVE-2026-33634",
      "date": "2026-03-26",
      "type": "Research"
    },
    {
      "name": "GitLab Security Advisory - CVE-2026-33634",
      "url": "https://advisories.gitlab.com/pkg/golang/github.com/aquasecurity/trivy/CVE-2026-33634/",
      "date": "2026-03-19",
      "type": "Vendor Advisory"
    },
    {
      "name": "Arctic Wolf - TeamPCP Supply Chain Attack Campaign",
      "url": "https://arcticwolf.com/resources/blog/teampcp-supply-chain-attack-campaign-targets-trivy-checkmarx-kics-and-litellm-potential-downstream-impact-to-additional-projects/",
      "date": "2026-03-27",
      "type": "Research"
    },
    {
      "name": "Help Net Security - CVE-2026-33634 Exploited in the Wild",
      "url": "https://www.helpnetsecurity.com/2026/03/27/cve-2026-33017-cve-2026-33634-exploited/",
      "date": "2026-03-27",
      "type": "News"
    },
    {
      "name": "PenLigent Security - Mutable Tags and the Trivy Supply Chain Compromise",
      "url": "https://www.penligent.ai/hackinglabs/cve-2026-33634-and-the-trivy-supply-chain-compromise-how-mutable-tags-turned-a-security-scanner-into-a-credential-stealer/",
      "date": "2026-03-21",
      "type": "Research"
    },
    {
      "name": "GitHub - Trivy Discussion #10425 (Official Incident Thread)",
      "url": "https://github.com/aquasecurity/trivy/discussions/10425",
      "date": "2026-03-19",
      "type": "Vendor Advisory"
    },
    {
      "name": "GitHub Security Advisory GHSA-69fq-xp46-6x23",
      "url": "https://github.com/aquasecurity/trivy/security/advisories/GHSA-69fq-xp46-6x23",
      "date": "2026-03-19",
      "type": "Vendor Advisory"
    },
    {
      "name": "Aqua Security Blog - Trivy Supply Chain Attack: What You Need to Know",
      "url": "https://www.aquasec.com/blog/trivy-supply-chain-attack-what-you-need-to-know/",
      "date": "2026-03-20",
      "type": "Vendor Advisory"
    },
    {
      "name": "Aikido Security - TeamPCP Deploys CanisterWorm via npm Trivy Compromise",
      "url": "https://www.aikido.dev/blog/teampcp-deploys-worm-npm-trivy-compromise",
      "date": "2026-03-21",
      "type": "Research"
    },
    {
      "name": "CISA - Adds CVE-2026-33634 to Known Exploited Vulnerabilities Catalog",
      "url": "https://www.cisa.gov/news-events/alerts/2026/03/26/cisa-adds-one-known-exploited-vulnerability-catalog",
      "date": "2026-03-26",
      "type": "Government"
    },
    {
      "name": "Docker Blog - Trivy Supply Chain Compromise: What Docker Hub Users Should Know",
      "url": "https://www.docker.com/blog/trivy-supply-chain-compromise-what-docker-hub-users-should-know/",
      "date": "2026-03-20",
      "type": "Vendor Advisory"
    },
    {
      "name": "OffSeq Threat Radar - CVE-2026-33634 Embedded Malicious Code",
      "url": "https://radar.offseq.com/threat/cve-2026-33634-cwe-506-embedded-malicious-code-in--163a34d0",
      "date": "2026-03-20",
      "type": "Research"
    },
    {
      "name": "CVE Details - CVE-2026-33634",
      "url": "https://www.cvedetails.com/cve/CVE-2026-33634/",
      "date": "2026-03-19",
      "type": "Government"
    },
    {
      "name": "Windows Forum - CISA Adds Trivy CVE-2026-33634 to KEV",
      "url": "https://windowsforum.com/threads/cisa-adds-trivy-cve-2026-33634-to-kev-patch-supply-chain-risk-now.407639/",
      "date": "2026-03-26",
      "type": "News"
    },
    {
      "name": "RedPacket Security - CVE-2026-33634 Alert",
      "url": "https://www.redpacketsecurity.com/cve-alert-cve-2026-33634-aquasecurity-setup-trivy/",
      "date": "2026-03-26",
      "type": "News"
    }
  ],
  "timeline": [
    {
      "date": "2026-02-01",
      "event": "Initial breach: Attackers exploit GitHub Actions misconfiguration in Aqua Security's Trivy repository, extracting a privileged access token"
    },
    {
      "date": "2026-03-01",
      "event": "Aqua Security discovers the February breach and begins credential rotation — but rotation is non-atomic, leaving residual access for attackers"
    },
    {
      "date": "2026-03-19",
      "event": "TeamPCP exploits retained credentials to force-push malicious commits to 76 of 77 trivy-action version tags and all 7 setup-trivy tags; publishes malicious Trivy v0.69.4 to all distribution channels (GitHub Releases, Docker Hub, GHCR, Amazon ECR Public, Debian/RPM repos)"
    },
    {
      "date": "2026-03-19",
      "event": "Malicious Trivy v0.69.4 binary and images detected and removed from registries (~4 hour window for binary, ~12 hours for GitHub Actions)"
    },
    {
      "date": "2026-03-21",
      "event": "CanisterWorm discovered propagating to 141+ npm packages using ICP blockchain canister as C2 infrastructure — first documented case of blockchain-based malware C2"
    },
    {
      "date": "2026-03-26",
      "event": "CISA adds CVE-2026-33634 to Known Exploited Vulnerabilities (KEV) catalog with federal remediation deadline of April 8, 2026"
    },
    {
      "date": "2026-03-27",
      "event": "Arctic Wolf publishes research revealing broader campaign targeting Checkmarx KICS, LiteLLM, and additional open-source security projects"
    }
  ],
  "simple_summary": {
    "what_happened": "Aqua Security's Trivy — the world's most widely used open-source container vulnerability scanner — was secretly compromised by a hacking group called TeamPCP. The attackers used stolen administrator credentials to inject credential-stealing malware into Trivy's official GitHub releases, Docker images, and automated CI/CD actions (used in thousands of software pipelines). Any organization that ran Trivy between March 19–20, 2026 may have had all their cloud credentials, API keys, SSH keys, and other secrets silently stolen. The malware also spread to 141+ npm software packages using blockchain technology for its command-and-control infrastructure — a first of its kind technique.",
    "who_is_affected": "Any organization using Trivy for container or infrastructure security scanning is potentially affected, particularly those using: Trivy v0.69.4 (any distribution channel), aquasecurity/trivy-action GitHub Action (versions 0.0.1 through 0.34.2), aquasecurity/setup-trivy GitHub Action (any version before 0.2.6), or any npm package infected by CanisterWorm (141+ packages). This affects cloud-native teams, DevOps engineers, and security teams across thousands of organizations globally.",
    "what_attackers_do": "Step 1: Attackers (TeamPCP) stole privileged GitHub credentials from Aqua Security in February 2026 by exploiting a GitHub Actions misconfiguration. Step 2: They waited while Aqua Security did an incomplete credential rotation, retaining backdoor access. Step 3: On March 19, they used the stolen credentials to replace the official Trivy binary and GitHub Actions with malware-infected versions across all distribution channels simultaneously. Step 4: When victims' CI/CD pipelines ran Trivy scans, the malware secretly harvested all accessible secrets — cloud credentials, SSH keys, API tokens, database passwords — from pipeline memory and the filesystem. Step 5: Stolen secrets were encrypted and exfiltrated to attacker servers, with a fallback of creating a hidden GitHub repository ('tpcp-docs') on the victim's account. Step 6: On developer machines, the malware installed a persistent backdoor that polls a blockchain smart contract every 5 minutes for new attack commands. Step 7: Using stolen npm tokens, CanisterWorm propagated automatically to 141+ npm packages, infecting downstream users who installed those packages.",
    "what_you_should_do": [
      "Immediately check if your organization used Trivy v0.69.4, trivy-action (any version before 0.35.0), or setup-trivy (any version before 0.2.6) on or around March 19–20, 2026",
      "Treat ALL secrets accessible to any affected CI/CD pipeline as fully compromised — rotate GitHub tokens, AWS/GCP/Azure credentials, SSH keys, Kubernetes service account tokens, Docker registry credentials, and database passwords simultaneously",
      "Check GitHub Actions workflow logs from March 19–20 for unexpected processes, environment variable dumps, or repository creation events (look for 'tpcp-docs' repository pattern)",
      "Update to safe versions: Trivy v0.69.3 or v0.69.5+, trivy-action v0.35.0+, setup-trivy v0.2.6+",
      "Pin all GitHub Actions to immutable commit SHAs rather than mutable version tags to prevent future tag-hijacking attacks"
    ],
    "why_this_is_serious": "This attack is exceptionally serious for three reasons. First, Trivy is trusted by thousands of organizations specifically to find security vulnerabilities — attackers weaponized that trust to steal the keys to cloud infrastructure. Second, the CanisterWorm's use of blockchain (ICP canister) for command-and-control infrastructure is unprecedented and cannot be taken down through conventional hosting provider takedowns. Third, the attack had a cascading effect: stolen npm tokens were used to propagate the malware to 141+ additional packages, meaning organizations that never used Trivy directly may still be compromised through transitive dependencies. The CISA KEV designation confirms active exploitation."
  },
  "mitre_attack": {
    "tactics": [
      {
        "tactic": "Initial Access",
        "technique_id": "T1195.002",
        "technique_name": "Supply Chain Compromise: Compromise Software Supply Chain",
        "description": "TeamPCP compromised Trivy's official distribution channels (GitHub Releases, Docker Hub, GHCR, Amazon ECR, Debian/RPM repos) by using stolen credentials to inject malicious code into legitimate builds"
      },
      {
        "tactic": "Execution",
        "technique_id": "T1053.006",
        "technique_name": "Scheduled Task/Job: Systemd Timers",
        "description": "On developer machines, a persistent systemd user service was installed to execute malicious payloads delivered via ICP canister on a 5-minute polling interval"
      },
      {
        "tactic": "Persistence",
        "technique_id": "T1547.001",
        "technique_name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
        "description": "Systemd user service (sysmon.py) installed for persistent execution, surviving reboots and operating invisibly within normal system service framework"
      },
      {
        "tactic": "Credential Access",
        "technique_id": "T1003",
        "technique_name": "OS Credential Dumping",
        "description": "Malware dumped Runner.Worker and Runner.Listener process memory via /proc/<pid>/mem to extract CI/CD runtime secrets and tokens"
      },
      {
        "tactic": "Credential Access",
        "technique_id": "T1552.001",
        "technique_name": "Unsecured Credentials: Credentials In Files",
        "description": "Systematically scanned 50+ filesystem paths for SSH keys, cloud credentials, kubeconfig files, .git/config, and other stored secrets"
      },
      {
        "tactic": "Defense Evasion",
        "technique_id": "T1036",
        "technique_name": "Masquerading",
        "description": "Attacker spoofed legitimate GitHub commit authors (rauchg for actions/checkout, DmitriyLewen for aquasecurity) to make malicious commits appear authentic"
      },
      {
        "tactic": "Exfiltration",
        "technique_id": "T1041",
        "technique_name": "Exfiltration Over C2 Channel",
        "description": "Stolen credentials encrypted with AES-256 + RSA-4096 and exfiltrated via HTTPS to C2 domain (scan[.]aquasecurtiy[.]org), with fallback to GitHub repository creation ('tpcp-docs')"
      },
      {
        "tactic": "Impact",
        "technique_id": "T1531",
        "technique_name": "Account Access Removal",
        "description": "Destructive wiper variant detected systems with Asia/Tehran timezone or fa_IR locale and executed rm -rf / to destroy data — targeted destructive component alongside credential theft"
      }
    ]
  },
  "detection_guidance": {
    "summary": "Detection focuses on identifying whether your environment executed the compromised Trivy artifacts between March 19-20, 2026, monitoring for credential exfiltration indicators, and detecting CanisterWorm activity in npm dependencies. Primary signals include unexpected process memory access in CI/CD logs, outbound connections to the typosquatted C2 domain, and suspicious GitHub repository creation.",
    "log_sources": [
      {
        "source": "GitHub Actions Workflow Logs",
        "description": "Check workflow run logs from March 19-20, 2026 for any jobs using trivy-action or setup-trivy. Look for unexpected environment variable dumps, process listing, or repository creation events",
        "priority": "Critical"
      },
      {
        "source": "Network/Firewall Egress Logs",
        "description": "Scan for outbound HTTPS connections to 'scan.aquasecurtiy.org' (typosquatter domain) and any unusual ICP blockchain network traffic",
        "priority": "Critical"
      },
      {
        "source": "GitHub Audit Logs",
        "description": "Monitor for unexpected repository creation events — particularly any repository named 'tpcp-docs' — which indicates the malware's fallback exfiltration mechanism was triggered",
        "priority": "High"
      },
      {
        "source": "CI/CD Process Logs",
        "description": "Look for unusual process executions spawned from Trivy scan steps, /proc/<pid>/mem file access patterns, or filesystem scanning operations outside expected scan scope",
        "priority": "High"
      },
      {
        "source": "systemd Journal Logs",
        "description": "On developer machines, check for unexpected new systemd user services, particularly 'sysmon.py' or similar Python-based services created around March 19-21, 2026",
        "priority": "High"
      },
      {
        "source": "npm Audit Logs",
        "description": "Run npm audit on all projects to detect CanisterWorm-infected packages in dependency trees; compare package hashes against known-good versions",
        "priority": "Medium"
      }
    ],
    "what_to_look_for": [
      {
        "indicator": "C2 Domain: scan[.]aquasecurtiy[.]org",
        "description": "Typosquatted domain (note: 'aquasecurtiy' not 'aquasecurity') used as primary exfiltration endpoint for stolen credentials",
        "log_field": "dns.query.name='scan.aquasecurtiy.org' OR url.domain='scan.aquasecurtiy.org'",
        "severity": "Critical"
      },
      {
        "indicator": "GitHub Repository 'tpcp-docs' created on victim account",
        "description": "Malware creates a hidden repository named 'tpcp-docs' as a fallback data exfiltration channel when the primary C2 is unreachable",
        "log_field": "github.event_type='CreateRepository' AND github.repo.name='tpcp-docs'",
        "severity": "Critical"
      },
      {
        "indicator": "Trivy v0.69.4 execution in CI/CD",
        "description": "Any CI/CD job that executed Trivy version 0.69.4 between March 19-20, 2026 should be treated as compromised",
        "log_field": "process.name='trivy' AND process.args CONTAINS '0.69.4'",
        "severity": "Critical"
      },
      {
        "indicator": "trivy-action versions 0.0.1-0.34.2 used in GitHub Actions",
        "description": "Any GitHub Actions workflow referencing aquasecurity/trivy-action@v0.0.1 through v0.34.2 during the compromise window",
        "log_field": "github.workflow.step.uses MATCHES 'aquasecurity/trivy-action@v0\\.[0-2]'",
        "severity": "High"
      },
      {
        "indicator": "ICP canister network queries",
        "description": "Unusual blockchain network traffic to Internet Computer Protocol canister endpoints indicates CanisterWorm C2 polling activity",
        "log_field": "url.domain CONTAINS '.ic0.app' OR url.domain CONTAINS 'icp0.io' AND process.name='python'",
        "severity": "High"
      },
      {
        "indicator": "Suspicious systemd user service 'sysmon.py'",
        "description": "CanisterWorm installs a persistent Python-based systemd user service disguised as system monitoring",
        "log_field": "systemd.unit.name='sysmon' AND systemd.unit.path CONTAINS '.config/systemd/user/'",
        "severity": "High"
      }
    ],
    "known_malicious_ips": [
      "scan[.]aquasecurtiy[.]org (typosquatted C2 domain)"
    ]
  },
  "immediate_actions": [
    {
      "priority": 1,
      "action": "Audit Trivy Usage on March 19-20, 2026",
      "description": "Check GitHub Actions workflow logs, CI/CD pipeline histories, and container orchestration logs to determine if your organization ran Trivy v0.69.4, trivy-action (pre-v0.35.0), or setup-trivy (pre-v0.2.6) during the compromise window",
      "command": "grep -r 'trivy' .github/workflows/ | grep -E '(v0\\.69\\.4|0\\.3[0-4]\\.)'"
    },
    {
      "priority": 2,
      "action": "Rotate ALL CI/CD Secrets Simultaneously",
      "description": "If any exposure is confirmed, treat every secret accessible to affected pipelines as compromised. Rotate GitHub tokens, AWS/GCP/Azure credentials, SSH keys, Kubernetes tokens, Docker credentials, and database passwords in a single coordinated operation — non-atomic rotation (as Aqua Security did) allows attackers to retain access",
      "command": null
    },
    {
      "priority": 3,
      "action": "Update to Safe Trivy Versions",
      "description": "Update trivy-action to v0.35.0+, setup-trivy to v0.2.6+, and Trivy binary to v0.69.3 (stable) or v0.69.5+. Pin all GitHub Actions to immutable commit SHAs instead of version tags to prevent future tag-hijacking",
      "command": "# Update trivy-action in workflows:\n# uses: aquasecurity/trivy-action@v0.35.0\n# Pin to SHA: uses: aquasecurity/trivy-action@<full-commit-sha>"
    },
    {
      "priority": 4,
      "action": "Check for CanisterWorm in npm Dependencies",
      "description": "Run npm audit and dependency analysis to identify any of the 141+ npm packages infected by CanisterWorm. Compare installed package versions and checksums against known-good versions from before March 19, 2026",
      "command": "npm audit --audit-level=moderate && npm list --all | grep -f known-compromised-packages.txt"
    },
    {
      "priority": 5,
      "action": "Search for Persistent Backdoor on Developer Machines",
      "description": "On any developer machine that may have pulled malicious artifacts, check for the CanisterWorm systemd user service and remove it. Also check for unauthorized GitHub repository 'tpcp-docs' on organization accounts",
      "command": "systemctl --user list-units --type=service | grep -i sysmon; ls ~/.config/systemd/user/ | grep -v snap"
    }
  ]
}
