# public/data/briefings/TB-2026-04-05-002.json
{
  "briefing_id": "TB-2026-04-05-002",
  "title": "TrueConf Client Update Integrity Bypass CVE-2026-3502 (CISA KEV, Deadline April 16): Chinese-Nexus Operation TrueChaos Deploys Havoc RAT Against Southeast Asian Governments via Tampered Software Updates and DLL Sideloading",
  "generated_date": "2026-04-05",
  "status": "ACTIVE",
  "severity": {
    "level": "High",
    "cvss": 7.8,
    "cisa_kev": true,
    "cisa_deadline": "2026-04-16"
  },
  "cves": [
    "CVE-2026-3502"
  ],
  "tags": [
    "trueconf",
    "cisa-kev",
    "supply-chain",
    "dll-sideloading",
    "havoc",
    "government",
    "espionage",
    "southeast-asia"
  ],
  "affected_products": [
    "TrueConf Client for Windows versions 8.1.0 through 8.5.2"
  ],
  "sources_analyzed": 15,
  "sources": [
    {
      "name": "CISA Known Exploited Vulnerabilities Catalog - CVE-2026-3502",
      "url": "https://www.cisa.gov/known-exploited-vulnerabilities-catalog",
      "date": "2026-04-02",
      "type": "Government"
    },
    {
      "name": "Check Point Research - Operation TrueChaos Technical Analysis",
      "url": "https://blog.checkpoint.com/research/when-trusted-software-updates-become-the-attack-vector-inside-operation-truechaos-and-a-new-zero-day-vulnerability-in-a-popular-collaboration-tool/",
      "date": "2026-03-31",
      "type": "Research"
    },
    {
      "name": "Check Point Research - Operation TrueChaos Zero-Day Exploitation",
      "url": "https://research.checkpoint.com/2026/operation-truechaos-0-day-exploitation-against-southeast-asian-government-targets/",
      "date": "2026-03-31",
      "type": "Research"
    },
    {
      "name": "GB Hackers - TrueConf Vulnerability Under Active Exploitation",
      "url": "https://gbhackers.com/trueconf-vulnerability-under-active-exploitation/",
      "date": "2026-04-02",
      "type": "Research"
    },
    {
      "name": "TheHackerWire - CVE-2026-3502 Vulnerability Details",
      "url": "https://www.thehackerwire.com/vulnerability/CVE-2026-3502/",
      "date": "2026-04-02",
      "type": "Research"
    },
    {
      "name": "THREATINT - CVE-2026-3502 Entry",
      "url": "https://cve.threatint.eu/CVE/CVE-2026-3502",
      "date": "2026-04-02",
      "type": "Research"
    },
    {
      "name": "TrueConf Official Blog - Version 8.5.3 Security Update",
      "url": "https://trueconf.com/blog/update/trueconf-8-5-3",
      "date": "2026-03-30",
      "type": "Vendor Advisory"
    },
    {
      "name": "TrueConf Official Changelog",
      "url": "https://trueconf.com/products/changelog.html",
      "date": "2026-03-30",
      "type": "Vendor Advisory"
    },
    {
      "name": "The Hacker News - TrueConf Zero-Day Exploited in Government Attacks",
      "url": "https://thehackernews.com/2026/03/trueconf-zero-day-exploited-in-attacks.html",
      "date": "2026-04-01",
      "type": "News"
    },
    {
      "name": "BleepingComputer - Hackers Exploit TrueConf Zero-Day to Push Malicious Software Updates",
      "url": "https://www.bleepingcomputer.com/news/security/hackers-exploit-trueconf-zero-day-to-push-malicious-software-updates/",
      "date": "2026-04-02",
      "type": "News"
    },
    {
      "name": "SecurityWeek - TrueConf Zero-Day Exploited in Asian Government Attacks",
      "url": "https://www.securityweek.com/trueconf-zero-day-exploited-in-asian-government-attacks/",
      "date": "2026-04-01",
      "type": "News"
    },
    {
      "name": "The Record (Recorded Future) - TrueConf Cyberattack CISA Notification",
      "url": "https://therecord.media/trueconf-cyberattack-cisa-hackers",
      "date": "2026-04-02",
      "type": "News"
    },
    {
      "name": "Security Affairs - CISA Adds TrueConf Client Flaw to Known Exploited Vulnerabilities Catalog",
      "url": "https://securityaffairs.com/190341/security/u-s-cisa-adds-a-flaw-in-trueconf-client-to-its-known-exploited-vulnerabilities-catalog.html",
      "date": "2026-04-02",
      "type": "News"
    },
    {
      "name": "Help Net Security - TrueConf Zero-Day Vulnerability Used in Cyber Espionage",
      "url": "https://www.helpnetsecurity.com/2026/04/02/trueconf-zero-day-vulnerability-cyber-espionage/",
      "date": "2026-04-02",
      "type": "News"
    },
    {
      "name": "InfoSec Today - CISA Adds TrueConf Flaw to Known Exploited Vulnerabilities",
      "url": "https://www.infosectoday.io/u-s-cisa-adds-a-flaw-in-trueconf-client-to-its-known-exploited-vulnerabilities-catalog/",
      "date": "2026-04-02",
      "type": "News"
    }
  ],
  "timeline": [
    {
      "date": "2026-01-01",
      "event": "Operation TrueChaos begins — Chinese-nexus threat actor starts targeting Southeast Asian government entities using compromised TrueConf update infrastructure"
    },
    {
      "date": "2026-03-03",
      "event": "CVE-2026-3502 reserved in the CVE database"
    },
    {
      "date": "2026-03-30",
      "event": "TrueConf releases version 8.5.3 containing the security fix for CVE-2026-3502; vulnerability publicly disclosed"
    },
    {
      "date": "2026-03-31",
      "event": "Check Point Research publishes full technical analysis of Operation TrueChaos including attack chain, IOCs, and threat attribution"
    },
    {
      "date": "2026-04-02",
      "event": "CISA adds CVE-2026-3502 to the Known Exploited Vulnerabilities (KEV) catalog; federal agency remediation deadline set for April 16, 2026"
    },
    {
      "date": "2026-04-16",
      "event": "Federal civilian executive branch agency mandatory remediation deadline per BOD 22-01"
    }
  ],
  "simple_summary": {
    "what_happened": "A Chinese-linked hacking group (tracked as Operation TrueChaos by Check Point Research) exploited a critical flaw in TrueConf — a popular video conferencing platform widely used by governments, defense organizations, and enterprises, particularly in Russia and Southeast Asia. The vulnerability (CVE-2026-3502) exists because TrueConf Client does not verify the authenticity of software updates it receives from the central server. Attackers who could compromise the TrueConf server were able to push malicious updates that installed a remote-access trojan (the Havoc post-exploitation framework) on every connected government client. CISA added this vulnerability to its mandatory remediation list on April 2, 2026, with a 14-day deadline for all US federal agencies.",
    "who_is_affected": "Organizations using TrueConf Client for Windows versions 8.1.0 through 8.5.2, particularly government agencies, defense organizations, and critical infrastructure that use TrueConf for internal video conferencing. The campaign has confirmed victims in Southeast Asian government entities. Any organization using on-premises TrueConf deployment where the central server could be accessible to threat actors is at risk. TrueConf version 8.5.3 and later contain the fix.",
    "what_attackers_do": "Step 1: Attackers compromise an organization's on-premises TrueConf server by exploiting other vulnerabilities, using stolen credentials, or through VPN compromise — this gives them control over the update delivery mechanism. Step 2: Attackers craft a malicious update package that appears legitimate, containing a benign-looking executable (poweriso.exe) bundled with a malicious DLL file (7z-x64.dll) designed to execute via DLL sideloading. Step 3: Attackers replace the legitimate update files on the TrueConf server with the malicious package, and update version metadata to indicate a new version is available. Step 4: When TrueConf clients check for updates, they receive a prompt saying a new version (8.5.3) is ready — users accept and trigger the download. Step 5: The TrueConf client downloads and installs the malicious update without performing any cryptographic verification or integrity checking (the core vulnerability). Step 6: During installation, poweriso.exe is launched and automatically loads the malicious 7z-x64.dll via DLL sideloading. Step 7: The malicious DLL performs reconnaissance (process listing, network tracing), bypasses Windows UAC via iscicpl.exe, and establishes registry-based persistence. Step 8: The Havoc post-exploitation framework is deployed, establishing persistent command-and-control communications — allowing attackers to spy on government networks, steal classified documents, and move laterally to other systems.",
    "what_you_should_do": [
      "Update all TrueConf Client installations to version 8.5.3 or later immediately — this is mandatory for US federal agencies by April 16, 2026",
      "Hunt for indicators of compromise on all systems with TrueConf installed: check for C:\\ProgramData\\PowerISO\\ directory, poweriso.exe, 7z-x64.dll, and iscsiexe.dll files",
      "Check the Windows registry for unauthorized run keys: HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\UpdateCheck pointing to C:\\ProgramData\\PowerISO\\PowerISO.exe",
      "Search network logs for outbound connections to known Havoc C2 infrastructure: IP 47.237.15[.]197 and related Alibaba Cloud / Tencent Cloud-hosted servers",
      "Immediately audit all on-premises TrueConf servers for unauthorized access, changed files, and any modifications to update packages or version metadata"
    ],
    "why_this_is_serious": "This is a nation-state level supply chain attack targeting government infrastructure — not opportunistic cybercrime. Every government agency that trusted its TrueConf server potentially received the backdoor as a routine software update, with no warning. The Havoc framework provides attackers with full control over compromised systems: keylogging, screen capture, credential theft, and lateral movement across government networks. The confirmed targeting of Southeast Asian governments suggests strategic intelligence-gathering with geopolitical implications. The fact that CISA issued a mandatory remediation deadline underscores the severity and active exploitation status."
  },
  "mitre_attack": {
    "tactics": [
      {
        "tactic": "Initial Access",
        "technique_id": "T1195.003",
        "technique_name": "Supply Chain Compromise: Compromise Software Supply Chain",
        "description": "Attackers compromised the on-premises TrueConf server (trusted software update infrastructure) to inject malicious code into the update delivery pipeline, exploiting CVE-2026-3502's lack of integrity verification"
      },
      {
        "tactic": "Persistence",
        "technique_id": "T1547.001",
        "technique_name": "Boot or Logon Autostart Execution: Registry Run Keys",
        "description": "Malware creates registry run key HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\UpdateCheck pointing to C:\\ProgramData\\PowerISO\\PowerISO.exe to ensure execution on every user logon"
      },
      {
        "tactic": "Privilege Escalation",
        "technique_id": "T1548.002",
        "technique_name": "Abuse Elevation Control Mechanism: Bypass User Account Control",
        "description": "Malicious DLL bypasses Windows User Account Control (UAC) using iscicpl.exe to elevate from standard user to administrator privileges without triggering a UAC prompt"
      },
      {
        "tactic": "Defense Evasion",
        "technique_id": "T1574.001",
        "technique_name": "Hijack Execution Flow: DLL Search Order Hijacking",
        "description": "Malicious 7z-x64.dll is loaded by the benign-appearing poweriso.exe executable via DLL sideloading, hiding malicious activity behind a trusted-looking process"
      },
      {
        "tactic": "Discovery",
        "technique_id": "T1057",
        "technique_name": "Process Discovery",
        "description": "Malware executes 'tasklist' to enumerate running processes during the reconnaissance phase, profiling the victim's security tools and environment"
      },
      {
        "tactic": "Command and Control",
        "technique_id": "T1071.001",
        "technique_name": "Application Layer Protocol: Web Protocols",
        "description": "Havoc post-exploitation framework establishes persistent HTTPS command-and-control communications with attacker-controlled servers hosted on Alibaba Cloud and Tencent Cloud infrastructure"
      },
      {
        "tactic": "Exfiltration",
        "technique_id": "T1041",
        "technique_name": "Exfiltration Over C2 Channel",
        "description": "Classified government documents, credentials, and sensitive data are exfiltrated through the established Havoc C2 channel"
      }
    ]
  },
  "detection_guidance": {
    "summary": "Detection should prioritize hunting for the known IOC file paths and registry keys planted by Operation TrueChaos. The presence of C:\\ProgramData\\PowerISO\\ directory on a system that has never had PowerISO legitimately installed is a high-fidelity indicator. Network detection should focus on outbound connections from TrueConf processes to non-TrueConf infrastructure, particularly known C2 IP 47.237.15[.]197. TrueConf update integrity monitoring should be implemented to catch future tampering.",
    "log_sources": [
      {
        "source": "Windows Registry (via EDR or Sysmon Event ID 13)",
        "description": "Monitor for creation of HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\UpdateCheck — a specific persistence mechanism planted by the Operation TrueChaos malware",
        "priority": "Critical"
      },
      {
        "source": "Windows File System (Sysmon Event ID 11)",
        "description": "Monitor for creation of C:\\ProgramData\\PowerISO\\poweriso.exe, 7z-x64.dll, and iscsiexe.dll — these are specific malware artifacts from the Operation TrueChaos campaign",
        "priority": "Critical"
      },
      {
        "source": "Network Traffic Logs / Firewall / Proxy",
        "description": "Monitor outbound connections from TrueConf client processes to external IPs not associated with TrueConf infrastructure; alert on connections to 47.237.15[.]197 (confirmed Havoc C2)",
        "priority": "Critical"
      },
      {
        "source": "Windows Security Event Logs (Event ID 4688 / Sysmon Event ID 1)",
        "description": "Process creation events showing tasklist.exe or tracert.exe spawned as children of poweriso.exe or TrueConf update processes are a strong behavioral indicator of malware reconnaissance",
        "priority": "High"
      },
      {
        "source": "TrueConf Server Audit Logs",
        "description": "Review TrueConf server-side logs for unauthorized access, file modification events on update packages, and changes to version metadata — server compromise is a prerequisite for this attack",
        "priority": "High"
      }
    ],
    "what_to_look_for": [
      {
        "indicator": "Malware artifact directory on endpoint",
        "description": "Presence of C:\\ProgramData\\PowerISO\\ directory on systems without legitimate PowerISO software installed is a definitive indicator of Operation TrueChaos compromise",
        "log_field": "TargetFilename LIKE 'C:\\ProgramData\\PowerISO\\%'",
        "severity": "Critical"
      },
      {
        "indicator": "Persistence registry run key",
        "description": "Registry run key UpdateCheck under HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run pointing to C:\\ProgramData\\PowerISO\\PowerISO.exe is a specific IOC for this campaign",
        "log_field": "TargetObject='HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\UpdateCheck'",
        "severity": "Critical"
      },
      {
        "indicator": "Known Havoc C2 infrastructure connection",
        "description": "Outbound network connection to IP 47.237.15.197 (confirmed Havoc C2 server used by Operation TrueChaos) indicates active post-exploitation activity",
        "log_field": "dst_ip='47.237.15.197'",
        "severity": "Critical"
      },
      {
        "indicator": "Reconnaissance commands spawned from PowerISO context",
        "description": "tasklist.exe or tracert.exe executing as child processes of poweriso.exe indicates the malicious DLL performing its reconnaissance phase",
        "log_field": "ParentImage LIKE '%poweriso.exe' AND Image IN ('tasklist.exe','tracert.exe','net.exe')",
        "severity": "High"
      }
    ],
    "known_malicious_ips": [
      "47.237.15.197"
    ]
  },
  "immediate_actions": [
    {
      "priority": 1,
      "action": "Update TrueConf Client to version 8.5.3 or later on all endpoints",
      "description": "This is mandatory for US federal agencies by April 16, 2026 (CISA KEV deadline). Download update packages directly from official TrueConf channels — do NOT use automated update from potentially compromised on-premises servers until server integrity is verified.",
      "command": null
    },
    {
      "priority": 2,
      "action": "Hunt for Operation TrueChaos IOCs across the environment",
      "description": "Immediately search all endpoints with TrueConf installed for the presence of C:\\ProgramData\\PowerISO\\ directory, associated DLL files, and the UpdateCheck registry run key. Any positive hit indicates active compromise.",
      "command": "reg query HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v UpdateCheck"
    },
    {
      "priority": 3,
      "action": "Block known Havoc C2 infrastructure at the network perimeter",
      "description": "Block outbound connections to IP 47.237.15.197 at firewall and proxy level. Also consider blocking outbound connections to Alibaba Cloud and Tencent Cloud IP ranges if not required for business operations.",
      "command": null
    },
    {
      "priority": 4,
      "action": "Audit on-premises TrueConf server integrity",
      "description": "Review TrueConf server access logs, file integrity on update packages (compare checksums against vendor-published hashes), and all administrative account activity. Server compromise is the prerequisite for this attack — if your server was compromised, assume all TrueConf clients that updated in 2026 may be backdoored.",
      "command": null
    },
    {
      "priority": 5,
      "action": "Disable automatic TrueConf updates until server integrity is confirmed",
      "description": "Prevent TrueConf clients from automatically pulling updates from the on-premises server until server integrity has been fully verified. Distribute the patched version 8.5.3 through a separate, verified channel.",
      "command": null
    }
  ]
}
